# Library/LaunchAgents/me.00dani.listenbrainz-mpd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>me.00dani.listenbrainz-mpd</string>
	<key>Program</key>
	<string>/Users/dani/.local/share/cargo/bin/listenbrainz-mpd</string>
	<key>RunAtLoad</key>
	<true/>
</dict>
</plist>
