# Library/LaunchAgents/me.00dani.mpd-now-playable.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>MPD_NOW_PLAYABLE_CACHE</key>
		<string>redis://localhost:6380/0?namespace=mpd-now-playable</string>
	</dict>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>me.00dani.mpd-now-playable</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/dani/.local/bin/mpd-now-playable</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
</dict>
</plist>
